# scratch/lims-backend-8091.service
[Unit]
Description=LIMS Go Backend Service (Port 8091)
After=network.target postgresql.service

[Service]
Type=simple
User=lims
Group=lims
WorkingDirectory=/home/lims/lims2/backend
ExecStart=/home/lims/lims2/backend/main
Restart=always
RestartSec=5
LimitNOFILE=65535
EnvironmentFile=/home/lims/lims2/backend/.env

[Install]
WantedBy=multi-user.target
